# pnl.crontab
# minute (0-59),
#    hour (0-23),
#       day of the month (1-31),
#          month of the year (1-12),
#             day of the week (0-6, 0=Sunday),
#       

00 02 * * 6 /projects/schiz/log/weekly.sh
00 03 * * * /projects/schiz/log/daily.sh
00 13 * * 1 /projects/schiz/log/log_intrustpdf.sh; /projects/schiz/log/mailintrust.sh reckbo mjc
